# ansible/roles/machine/files/etc/update-motd.d/fortune
#!/bin/sh

# shellcheck disable=SC1091
. /etc/default/locale

export LANG
export PATH="/usr/local/games:/usr/games:$PATH"

if [ -x /usr/games/fortune ] ; then
    /usr/games/fortune -s
fi
