# templates/etc/systemd/system/elasticsearch.service.d/ansible.conf.j2
# {{ ansible_managed }}

[Service]
LimitMEMLOCK={{ elasticsearch__systemd_limit_memlock }}
