# templates/etc/systemd/system/home-assistant.service.j2
[Unit]
Description=Home Assistant
After=network.target

[Service]
Type=simple
User={{ homeassistant__user }}
ExecStart={{ homeassistant__prog_path }}

[Install]
WantedBy=multi-user.target
