# playbooks/roles/ginas.postgresql/files/etc/cron.daily/autopostgresqlbackup
#!/bin/sh

if [ -x /usr/sbin/autopostgresqlbackup ]; then
        /usr/sbin/autopostgresqlbackup
fi
