# ansible/roles/proc_hidepid/templates/etc/systemd/system/systemd-logind.service.d/hidepid.conf.j2
# {{ ansible_managed }}

# Based on https://github.com/indiv0/ansible-role-console

[Service]
SupplementaryGroups={{ proc_hidepid__group }}
