# ansible/roles/debops.elasticsearch/templates/etc/systemd/system/elasticsearch.service.d/ansible.conf.j2
# {{ ansible_managed }}

[Service]
{% if elasticsearch__memory_lock|bool %}
LimitMEMLOCK={{ elasticsearch__systemd_limit_memlock }}
{% endif %}
