# ansible/roles/resolvconf/templates/etc/systemd/system/resolvconf.service.d/static.conf.j2
# {{ ansible_managed }}

[Service]
ExecStartPost=-/usr/local/lib/resolvconf-static
