# ansible/roles/influxdb_server/templates/etc/cron.daily/autoinfluxdbbackup.j2
#!/bin/sh

# {{ ansible_managed }}

test -x /usr/sbin/autoinfluxdbbackup || exit 0
/usr/sbin/autoinfluxdbbackup > /dev/null
