# ansible/roles/debops.redis_server/templates/etc/systemd/system/redis-server.service.j2
# {{ ansible_managed }}

[Unit]
Description=Advanced key-value store
After=network.target
Documentation=http://redis.io/documentation, man:redis-server(1)

[Service]
Type=oneshot
ExecStart=/bin/true
ExecReload=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
Alias=redis.service
